# root/etc/systemd/system/rc.service
[Unit]
Description=Service that runs all commands found inside /etc/rc

[Service]
ExecStart=/etc/rc

[Install]
WantedBy=multi-user.target
